# unit/config-drive.path
[Unit]
Description=Config Drive unit for /media/configdrive

[Path]
DirectoryNotEmpty=/media/configdrive
