# unit/config-drive.mount
[Unit]
Wants=config-drive.path
Before=config-drive.service

[Mount]
Where=/media/configdrive
Options=ro
