# unit/flamingo.service
[Unit]
Description=Flamingo Cloud Instance Contextualization Tool

[Service]
ExecStart=-/bin/flamingo

[Install]
WantedBy=multi-user.target
