# pkg/sys/initd/systemd/test_units/echo.socket
[Unit]
Description=echo socket

[Socket]
ListenStream=2222
Accept=yes
