# pkg/sys/initd/systemd/test_units/echo@.service
[Unit]
Description=echo service

[Service]
ExecStart=-/bin/cat
StandardInput=socket
